# templates/supervisor.conf.j2
[program:dmm-gacha]
command=java -jar {{ artifact }}
directory= {{ home_dir }}
user=root
autostart=true
autorestart=true
redirect_stderr=true
