# templates/supervisor.conf.j2
[program:dmm-gacha]
command=java -jar dmm-gacha.jar
directory= {{dmm_home}}
user=root
autostart=true
autorestart=true
redirect_stderr=true
